# tasks/module_1/8_System_initialization_Systemd_and_SysV/8_task/1_watchlog/usr/lib/systemd/system/watchlog.timer
[Unit]
Description=Run watchlog script every 2 second

[Timer]
OnBootSec=1
OnUnitActiveSec=2
AccuracySec=1us
Unit=watchlog.service

[Install]
WantedBy=multi-user.target
